# config/saorview-mux2.service
[Unit]
Description=Saorview DTT Mux 2
After=network-online.target
 
[Service]
ExecStart=/bin/bash /home/conor/src/saorview-dvb-config/bin/mux2.sh
WorkingDirectory=/home/conor/src/saorview-dvb-config/bin
StandardOutput=inherit
StandardError=inherit
Restart=always
User=conor
 
[Install]
WantedBy=multi-user.target
